# roles/webgateway/files/certbot.target
[Unit]
Description=Script based certificate renewal via certbot

[Install]
WantedBy=default.target