# roles/webgateway/files/certbot-script.service
[Unit]
Description=Renew certificates with certbot

[Service]
Type=simple
ExecStart=/usr/bin/bash /opt/certbot.sh

[Install]
WantedBy=certbot.target