# roles/backupscript/files/backup-script.service
[Unit]
Description=Backup application data

[Service]
Type=simple
ExecStart=/usr/bin/bash /opt/backupscript.sh

[Install]
WantedBy=backup.target